# m2.conf
description "uWSGI server instance configured to serve m2"

start on runlevel [2345]
stop on runlevel [!2345]

setuid root
setgid www-data

env PATH=~/tn/env3/bin
chdir ~/tn/
